# templates/vault.service.tmpl
[Unit]
Description=Vault

[Service]
Type=simple
ExecStart={{vault_path}} server {{vault_options}}
